# deploy/varena-gateway.box.service
[Unit]
Description=vArena gateway (Firebase auth + proxy)
After=network.target varena.service
Wants=varena.service

[Service]
Type=simple
User=vseadmin
WorkingDirectory=/home/vseadmin/varena/gateway
EnvironmentFile=/home/vseadmin/varena/gateway/.env
ExecStart=/usr/bin/node server.mjs
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
